# etc/systemd/system/OpenLPOSCControl.service
[Unit]
Description=OpenLP OSC Control Service

[Service]
ExecStart=/usr/local/bin/OpenLPOSCControl

[Install]
WantedBy=default.target
